# home/user/.config/systemd/user/kerbnet.timer
[Unit]
Description=kerbnet.sh

[Timer]
OnBootSec=10min
OnUnitActiveSec=60min
Unit=kerbnet.service

[Install]
WantedBy=timers.target
